# systemd/arcade-control.service
[Unit]
Description=Arcade Control Panel
After=local-fs.target network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/alberto/apps/upac
ExecStart=/bin/bash /home/alberto/apps/upac/start_x11.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# X11 via xinit launched from console (no existing display)
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
